# public/config/etc/profile.d/custom.sh
source /public/opt/intel/parallel_studio_xe_2016.1.056/psxevars.sh
source /public/config/env-load.sh

